# server/cloudflared-avatar.service
[Unit]
Description=AI Meow Avatar Cloudflare private tunnel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/jerry/bin/cloudflared tunnel --no-autoupdate run --token-file /home/jerry/.cloudflared/avatar-tunnel.token
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
